# scripts_CI/scripts/gitpullappointments.service
[Unit]
Description=Bahmni continuous deployment
After=docker.service
After=mysqlbootstrap.service
After=dockerbootstrap.service


[Service]
WorkingDirectory=/development/openmrs-module-appointments-release
ExecStart=/development/openmrs-module-appointments-release/gitpull_appointments.sh


User=openmrs
RestartSec=10
Restart=always


[Install]
WantedBy=default.target
